# etc/profile.d/01_common.sh
## self common setting

# ls colors setting
export LS_COLORS='di=01;35:ln=04'
