# etc/profile.d/09-editor.sh
### editor

export EDITOR=me
export VISUAL=me
