# etc/profile.d/09-editor.sh
### editor

export EDITOR=em
export VISUAL=em
