# etc/profile.d/09-editor.sh
### editor

export EDITOR=ec
export VISUAL=ec
